# crates/mercury/assets/hg.freddie.mercury.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>hg.freddie.mercury</string>

    <!-- `daemon`, never the bare binary: `mercury` spawns a detached daemon and exits, which
         would leave launchd watching the job vanish and unable to see the process that actually
         holds the keyboard. `daemon` is hidden from `mercury --help` and stays invocable for
         exactly this and for `mercury start`. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/mercury</string>
        <string>daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Aqua: the session `CGEventTap` needs a window server, `NSWorkspace`, and per-user TCC,
         none of which a root daemon at the login window has. -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <!-- Revive a mercury that died unexpectedly; leave one down that declined to run. Both halves
         are the exit code: every deliberate way out exits zero (`q` from home, the menu bar's
         Quit, `mercury stop`, `launchctl bootout`, all of which reach the model's quit), and so
         does every refusal to start (the lock already held, the menu bar unavailable, the keyboard
         grab denied). None of those is fixed by trying again. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- A crash loop must not respawn something that eats the keyboard ten times a second. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- No StandardOutPath and no LOG_LEVEL. mercury writes ~/Library/Logs/mercury/mercury.log
         itself, and `--log-level` governs a terminal a launchd job does not have. -->
</dict>
</plist>
